# debian/postinst
#!/bin/sh

set -e

package="musica"
webserver="apache2"

case "$1" in
	configure)
		chown www-data:www-data /var/lib/$package/tmp
		a2ensite $package
		invoke-rc.d $webserver reload
		echo "Rebuilding musica cache..."
		/etc/cron.hourly/musica
	;;
	abort-upgrade|abort-remove|abort-deconfigure)
	;;
	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
